# dockerfiles/nomad/files/etc/s6/dockerd/run
#!/bin/sh

exec /usr/local/bin/dockerd-entrypoint.sh
